# scripts/oci/vbox-fcos/template/config-pull.service
[Unit]
Description=Perform Ansible pull from repository 

[Service]
ExecStart=xargs -a /etc/ansible.d/ansibled-pull.cfg ansible-pull